# steps/data/setup/puppet/modules/env/files/big/nvidia/ganglia-monitor-3.7-override.service
# INSTALLED BY PUPPET

[Service]
ExecStartPre=/bin/bash -c "[[ $(lsmod | grep -ic nvidia) -eq 0 ]] && rm -f /etc/ganglia/conf.d/{nvidia.pyconf,modpython-nvidia.conf} || true"
