# steps/data/setup/puppet/modules/env/files/std/fix_swap/fix_swap.service
[Install]
WantedBy=swap.target

[Unit]
Description=Force swapon fix page size

[Service]
Type=oneshot
ExecStart=swapon -f -a
